# setup/gpu-server/systemd/llama-chat.service.template
[Unit]
Description=llama.cpp Chat Server (Qwen 3.5 9B)
After=network.target

[Service]
ExecStart=/home/GPUUSER/llama.cpp/build/bin/llama-server \
  --model /home/GPUUSER/models/Qwen3.5-9B-Opus-Distilled-v2-Q4_K_M.gguf \
  --port 8080 --host 0.0.0.0 \
  --ctx-size 98304 \
  --n-gpu-layers 99 \
  --parallel 2 \
  --cache-type-k q4_0 --cache-type-v q4_0 \
  --reasoning-budget 1024 \
  --chat-template-file /home/GPUUSER/llama.cpp/chat-template.jinja
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
